# .config/systemd/user/battery-status-check.service
[Unit]
Description=Toggle Blur My Shell extension based on battery status
After=graphical-session.target

[Service]
Type=oneshot
ExecStart=/bin/bash /usr/local/bin/check-battery-and-toggle-blur
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
